# deploy/systemd/laravel-queue.service
[Unit]
Description=Laravel Queue Worker for Eurostar (Train Rules & GTFS Downloads)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
Restart=always
RestartSec=3
WorkingDirectory=/var/www/Eurostar
ExecStart=/usr/bin/php /var/www/Eurostar/artisan queue:work --queue=default,train-rules --sleep=3 --tries=3 --max-time=3600 --memory=512
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 