# deploy/supervisor/laravel-scheduler.conf
[program:laravel-scheduler]
process_name=%(program_name)s
command=php /var/www/Eurostar/artisan schedule:work
directory=/var/www/Eurostar
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-scheduler.log
stopwaitsecs=10 